# systemd/.config/systemd/user/gpg-agent.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/d.h47uaz7wfp511xezfcts17ot/S.gpg-agent