# systemd/.config/systemd/user/gpg-agent-ssh.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=ListenStream=%t/gnupg/S.gpg-agent.ssh